# supervisor/client.ini
[program:client_deamon]
command=python bin/deamon.py
numprocs=1
directory=/data/raspberry-pi-wechat-client
autostart=true
autorestart=true
startretries=0
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
environment=PYTHONPATH='/data/raspberry-pi-wechat-client'
