# supervisor/client.ini
[program:client_deamon]
command=python bin/deamon.py
numprocs=1
directory=/data/raspberry-pi-wechat-client
pythonpath=/data/raspberry-pi-wechat-client
autostart=true
autorestart=true
startretries=0
exitcodes=0
stopsignal=KILL
stopwaitsecs=10
redirect_stderr=true
